# socat_init
#!/sbin/openrc-run

command="/bin/server_socat.sh"

#command_background=true
pidfile="/run/socat.pid"

